# ops/systemd/samantha-telegram-poll.service
[Unit]
Description=Samantha Codex Telegram poll adapter
Documentation=https://github.com/lbk0523/samantha-codex
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/projects/samantha-codex
EnvironmentFile=-%h/projects/samantha-codex/.env
ExecStart=%h/.bun/bin/bun run samantha telegram:poll --timeout-seconds=25
TimeoutStartSec=45
