# ops/systemd/samantha-telegram-poll.timer
[Unit]
Description=Run Samantha Codex Telegram poll adapter
Documentation=https://github.com/lbk0523/samantha-codex

[Timer]
OnBootSec=10s
OnUnitInactiveSec=3s
AccuracySec=1s
Unit=samantha-telegram-poll.service

[Install]
WantedBy=timers.target
