# ops/systemd/samantha-inbox-watch.service
[Unit]
Description=Samantha Codex inbox watcher
Documentation=https://github.com/lbk0523/samantha-codex

[Service]
Type=simple
WorkingDirectory=%h/projects/samantha-codex
EnvironmentFile=-%h/projects/samantha-codex/.env
ExecStart=%h/.bun/bin/bun run samantha inbox:watch --interval-ms=1000
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=15

[Install]
WantedBy=default.target
