# ops/systemd/samantha-telegram-reply.timer
[Unit]
Description=Run Samantha Codex Telegram outbox reply adapter
Documentation=https://github.com/lbk0523/samantha-codex

[Timer]
OnBootSec=45s
OnUnitActiveSec=30s
AccuracySec=5s
Unit=samantha-telegram-reply.service

[Install]
WantedBy=timers.target
